# redis_server/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf

[program:user_permission]
command=/bin/bash -c "chown -R redis:syslog /var/log/redis/"

