# database_server/supervisord.conf
[supervisord]
nodaemon=true
 
[program:mysql]
command=/usr/bin/mysqld_safe
 
[program:script_of_user_creation]
command=/bin/bash -c "chmod a+x /etc/mysql.sh; /etc/mysql.sh"
